# template/dhcpd/rc.local
# /etc/rc.d/rc.local
...
# Start DHCP server
if [ -x /etc/rc.d/rc.dhcpd ]; then
  /etc/rc.d/rc.dhcpd start
fi
